# SomfyProtect2MQTT/rootfs/etc/s6-overlay/s6-rc.d/SomfyProtect2MQTT/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: SomfyProtect2MQTT
# Runs SomfyProtect2MQTT
# ==============================================================================

bashio::log.info "Generating config.yaml from options.json"
echo '# Generated by homeassistant, do not edit!' > /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml.tmp
echo '# Edit configuration only at the Add-on configuration tab!' >> /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml.tmp
yq -P -oy /data/options.json >> /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml.tmp
sed 's/  -/    -/'  /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml.tmp > /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml
bashio::log.info "Configuration :"
bashio::log.info "$(cat /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml)"

debug=""
if $(bashio::config.true debug); then
    bashio::log.info "Will be started in debug mode"
    debug="-v"
fi

bashio::log.info "Image build"
bashio::log.info "starting original stuff..."
exec python3 /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/main.py -c /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml $debug
