# SomfyProtect2MQTT/rootfs/etc/s6-overlay/s6-rc.d/SomfyProtect2MQTT/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: SomfyProtect2MQTT
# Runs SomfyProtect2MQTT
# ==============================================================================

bashio::log.info "Generating config.yaml from options.json"
export TZ="$(bashio::supervisor.timezone)"
echo '# Generated by homeassistant, do not edit!' > /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml.tmp
echo '# Edit configuration only at the Add-on configuration tab!' >> /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml.tmp
yq -P -oy /data/options.json >> /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml.tmp
sed 's/  -/    -/'  /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml.tmp > /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml
# Log configuration with masked passwords
bashio::log.info "Configuration:"
bashio::log.info "$(sed 's/\(password:\s*\).*/\1****/' /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml)"

if ! bashio::config.has_value 'mqtt.host'; then
    bashio::log.info "mqtt.host not set, adding to config.yaml with MQTT service value."
    yq -i ".mqtt.host = \"$(bashio::services mqtt \"host\")\"" /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml
fi
if ! bashio::config.has_value 'mqtt.port'; then
    bashio::log.info "mqtt.port not set, adding to config.yaml with MQTT service value."
    yq -i ".mqtt.port = $(bashio::services mqtt \"port\")" /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml
fi
if ! bashio::config.has_value 'mqtt.username'; then
    bashio::log.info "mqtt.username not set, adding to config.yaml with MQTT service value."
    yq -i ".mqtt.username = \"$(bashio::services mqtt \"username\")\"" /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml
fi
if ! bashio::config.has_value 'mqtt.password'; then
    bashio::log.info "mqtt.password not set, adding to config.yaml with MQTT service value."
    yq -i ".mqtt.password = \"$(bashio::services mqtt \"password\")\"" /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml
fi
if ! bashio::config.has_value 'mqtt.ssl'; then
    bashio::log.info "mqtt.ssl not set, adding to config.yaml with MQTT service value."
    yq -i ".mqtt.ssl = $(bashio::services mqtt \"ssl\")" /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml
fi

# Log final configuration with masked passwords
bashio::log.info "Final configuration:"
bashio::log.info "$(sed 's/\(password:\s*\).*/\1****/' /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml)"

debug=""
if $(bashio::config.true debug); then
    bashio::log.info "Will be started in debug mode"
    debug="-v"
fi

bashio::log.info "Image build"
bashio::log.info "starting original stuff..."
exec python3 /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/main.py -c /usr/bin/SomfyProtect2MQTT/somfyProtect2Mqtt/config/config.yaml $debug
